# install_scripts/scheduled_reports.crontab
0 0 * * * monitor /usr/bin/php /opt/monitor/op5/ninja/index.php default/cron/Daily >/dev/null 2>&1
0 0 * * 1 monitor /usr/bin/php /opt/monitor/op5/ninja/index.php default/cron/Weekly >/dev/null 2>&1
0 0 1 * * monitor /usr/bin/php /opt/monitor/op5/ninja/index.php default/cron/Monthly >/dev/null 2>&1
